# scripts/os-clean-posts.timer
[Unit]
Description=Clean openstock posts daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
