# scripts/os-delete-invitations.timer
[Unit]
Description=Delete openstock old invitation daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
